# linux/deb/prerm
#!/bin/bash
systemctl stop      supernode
systemctl disable   supernode
systemctl stop      edgenode
systemctl disable   edgenode
systemctl daemon-reload
exit 0;